# stock-insight.service
[Unit]
Description=StockInsight - 股票洞察分析系统
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/stock-insight
ExecStart=/usr/bin/python3 /opt/stock-insight/start_prod.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量（可选）
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

